# debian/postinst
#!/bin/sh
# postinst script for sickrage
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package


case "$1" in
    configure)
    chown -R admin:admin /apps/sickrage
    mkhomedir_helper admin 027

    #Check/create config directory
    if [ ! -d "/home/admin/.sickrage" ]; then
      mkdir /home/admin/.sickrage
      chown -R admin:admin /home/admin/.sickrage
      fi
    event_push app readynasd '<add-s resource-type="LocalApp" resource-id="LocalApp"><LocalApp appname="sickrage" success="1"/></add-s>' 0 0
    ;;
    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

###
# Enable service
###
if [ -f "/lib/systemd/system/fvapp-sickrage.service" ]; then
  systemctl enable fvapp-sickrage.service
  systemctl start fvapp-sickrage.service
fi


# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
